# deploy/nyx-serve-web.service
[Unit]
Description=serve.nyxlang.com landing page
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/nyx-sites/serve.nyxlang.com
ExecStart=/home/admin/nyx-sites/serve.nyxlang.com/serve-nyxlang-com
Environment=PORT=3003
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
